# install/thinknagent.service
[Unit]
Description=ThinkNCollab Agent
Documentation=https://thinkncollab.com/docs/agent
After=network-online.target
Wants=network-online.target

[Service]
User=thinknagent
Group=thinknagent
ExecStart=/usr/bin/env thinknagent start
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=%h/.thinknagent

PrivateTmp=true
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal
SyslogIdentifier=thinknagent

[Install]
WantedBy=multi-user.target
